# home/pi/.config/autostart/piScreenCore.desktop
[Desktop Entry] 
Name=PiScreenCore
Type=Application
Exec=/home/pi/piScreen/piScreenCore.py
Terminal=false
