# home/pi/.config/autostart/piScreenCore.desktop
[Desktop Entry] 
Name=PiScreenCore
Type=Application
Exec=sudo systemctl restart piScreen
Terminal=false
